# app-i18n/ibus-chewing/latest/root/etc/profile.d/30-ibus.sh
#!/bin/bash

export XMODIFIERS="@im=ibus"
export GTK_IM_MODULE="ibus"
export QT_IM_MODULE="ibus"

ibus-daemon -x -d

